# templates/archfrican-resume.service
# Archfrican first-boot resume (Stage 2). Installed into the NEW system by
# lib/phase1.sh::inject_resume after archinstall lays down the base. It finishes
# what the base install can't: the niri desktop + dev layer (the same modules the
# booted Stage-1 path runs). The @USER@ placeholder is substituted at install time.
#
# Lifecycle: runs ONCE as the wheel user. phase 2 is non-interactive here
# (ARCHFRICAN_NONINTERACTIVE=1, no TTY) so it reads ~/.archfrican-answers instead
# of prompting. On SUCCESS the ExecStartPost lines touch a marker file (and best-effort
# disable this unit + remove the temporary NOPASSWD sudoers drop-in) — so the elevated
# window is exactly one boot. On failure it retries next boot (phase-2 .done checkpoints
# make that safe). The elevated window is BOUNDED: ExecStartPre (lib/resume-guard.sh)
# counts boots and, after ARCHFRICAN_RESUME_MAX_BOOTS (default 5), touches the same
# marker file. ConditionPathExists=! below is what actually stops future boots from
# starting this unit again — systemd evaluates it with its own privilege, so (unlike the
# old sudo-systemctl-disable-based approach) it can never be blocked by a broken NOPASSWD
# grant. ARCHFRICAN_STRICT_PREFLIGHT=1 makes an unresolvable package fail loudly at
# preflight BEFORE any state change, instead of mid-module.
[Unit]
Description=Archfrican first-boot install (niri desktop + dev layer)
After=network-online.target systemd-user-sessions.service
Wants=network-online.target
ConditionPathExists=/home/@USER@/.archfrican/install.sh
ConditionPathExists=!/home/@USER@/.local/state/archfrican/resume-stopped

[Service]
Type=oneshot
RemainAfterExit=yes
User=@USER@
WorkingDirectory=/home/@USER@/.archfrican
Environment=HOME=/home/@USER@
Environment=ARCHFRICAN_NONINTERACTIVE=1
# Make a dropped/renamed package fail loudly at preflight (before any module mutates state).
Environment=ARCHFRICAN_STRICT_PREFLIGHT=1
TimeoutStartSec=infinity
# Fail-safe: bound the elevated retry window — resume-guard.sh touches the stop marker after N boots.
ExecStartPre=/usr/bin/bash /home/@USER@/.archfrican/lib/resume-guard.sh
ExecStart=/usr/bin/bash /home/@USER@/.archfrican/install.sh
# The marker is what actually stops future boots (via ConditionPathExists=! above, no sudo needed).
# The two sudo lines below are best-effort cleanup only, not required for correctness.
ExecStartPost=/usr/bin/touch /home/@USER@/.local/state/archfrican/resume-stopped
ExecStartPost=/usr/bin/sudo /usr/bin/systemctl disable archfrican-resume.service
ExecStartPost=/usr/bin/sudo /usr/bin/rm -f /etc/sudoers.d/99-archfrican-resume

[Install]
WantedBy=multi-user.target
